# meta-inspur/meta-fp5280g2/recipes-inspur/inspur-ncsi-switch/inspur-ncsi-card-check/inspur-ncsi-card-check.service
[Unit]
Description=Inspur NCSI Card Check
After=mapper-wait@-inspur-ncsi-manager.service

[Service]
ExecStart=inspur-ncsi-card-check.sh

[Install]
WantedBy=multi-user.target
